# unix/dot_config/systemd/user/dotfiles-update-check.timer
[Unit]
Description=Periodic package update check.

[Timer]
OnBootSec=0m
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
