# unix/dot_config/systemd/user/dotfiles-update-check.service
[Unit]
Description=Check and count available package updates.

[Service]
Type=oneshot
ExecStart="%h/scripts/dependencies/update-check.sh"
